# resources/root/etc/profile.d/16-custom-variables.sh
# Please keep this file empty for development version.
# On configured system it should be filled with custom variables which
# most likely comes from deployment script.

